# deploy/bytebrew-cloud-api.service
[Unit]
Description=ByteBrew Cloud API
After=network.target postgresql.service

[Service]
Type=simple
User=bytebrew
Group=bytebrew
WorkingDirectory=/etc/bytebrew
ExecStart=/usr/local/bin/bytebrew-cloud-api --config /etc/bytebrew/config.yaml
Restart=on-failure
RestartSec=5s
TimeoutStopSec=15s
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
