# deploy/bytebrew-bridge.service
[Unit]
Description=ByteBrew Bridge (gRPC relay)
After=network.target

[Service]
Type=simple
User=bytebrew
Group=bytebrew
EnvironmentFile=/etc/bytebrew/bridge.env
ExecStart=/usr/local/bin/bytebrew-bridge
Restart=on-failure
RestartSec=5s
TimeoutStopSec=15s
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
